# container-files/etc/supervisor.d/memcached.conf
[program:memcached]
command = /usr/bin/memcached -u root
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE