# container-files/etc/supervisor.d/sendmail.conf
[program:sendmail]
command = /usr/sbin/sendmail
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE